# fisher.service
[Unit]
Description=autostart script

[Service]
ExecStart=/home/sam/autostart.sh

[Install]
WantedBy=multi-user.target
